# cron.d/webchecks
# cron.d/webchecks
* *	* * *	root	( /usr/bin/python3 /opt/aws-webchecks/webchecks.py >> /var/log/webchecks.log 2>&1 )

